# supervisord.haproxy-ddos.conf
[program:haproxy-ddos]
command=/etc/haproxy-ddos/haproxy.sh
autorestart=unexpected
startretries=10
stdout_logfile=/var/log/haproxy-ddos-stdout.log
redirect_stderr=true
stderr_logfile=/dev/null